# 50-sd-card.rules
ACTION=="add", KERNEL="sd[a-z][0-9], ENV{ID_FS_TYPE}=="exfat", RUN+="/bin/systemctl start sd-mount@%k.service"
